# modules/docker/data/com.docker.machine.default.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
This is a launch agent for a docker-machine.
It is based on https://gist.github.com/andystanton/257fab335b242bc2658b, but
with setting of the PATH environment variable modified to be compatible with
the Yosemite version of launchd.
-->
    <dict>
        <key>Label</key>
        <string>com.docker.machine.default</string>
        <key>ProgramArguments</key>
        <array>
            <string>bash</string>
            <string>-c</string>
            <string>
                export PATH=/usr/local/bin:/usr/bin:/bin
                /usr/local/bin/docker-machine start default
            </string>
        </array>
        <key>RunAtLoad</key>
        <true/>
    </dict>
</plist>

